# systemd-services/order-api.service
[Unit]
Description=Order Automation API Server
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/order-automation
Environment=PATH=/opt/order-automation/venv/bin
Environment=PYTHONPATH=/opt/order-automation
ExecStart=/opt/order-automation/venv/bin/uvicorn app.main_api:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 